# third-party/5factorial-monitoring.service
[Unit]
Description=Service managing back monitoring
Requires=docker.service
After=docker.service

[Service]
User=root
ExecStartPre=/usr/local/bin/docker-compose --file /home/ubuntu/backend/2019_1_5factorial-team/third-party/docker-compose_monitoring.yml pull
ExecStart=/usr/local/bin/docker-compose --file /home/ubuntu/backend/2019_1_5factorial-team/third-party/docker-compose_monitoring.yml up
ExecStop=/usr/local/bin/docker-compose --file /home/ubuntu/backend/2019_1_5factorial-team/third-party/docker-compose_monitoring.yml down -v

[Install]
WantedBy=multi-user.target